# dist/etc/systemd/system/hermes.service
[Unit]
Description=daemon process for the hermes gateway and load balancer


[Service]
Type=simple
ExecStart=/opt/hermes/hermes
Restart=on-failure

[Install]
WantedBy=multi-user.target